# FreeBSD/bin/rc.local
#!/bin/sh

_NO_WRITE_STDERR=0 _bsdinstall &

exit 0
